# srcs/InfluxDB/srcs/supervisor.conf
[supervisord]
nodaemon=true
user = root

[program:influxdb]
command= sh -c "influxd run -config /etc/influxdb.conf && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command= sh -c "/telegraf/usr/bin/telegraf --config etc/telegraf/telegraf.conf \
				&& kill -s SIGTERM $(cat supervisord.pid)"